# deploy/basketball-api.service.example
# 复制为 /etc/systemd/system/basketball-api.service 后：
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now basketball-api
#
# 将 User/Group 改为能读取 /data/node_app/basketball_drill 的部署用户（与目录属主一致）。

[Unit]
Description=Basketball drill API (Node)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/data/node_app/basketball_drill
EnvironmentFile=/data/node_app/basketball_drill/.env
ExecStart=/usr/bin/node apps/api/dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
